# atlas-control.service
# NOTE: This file is a reference template only.
# The actual service file is generated by install.sh with the correct
# user, home directory, and app path substituted at install time.
# Do not copy this file directly to /etc/systemd/system/.

[Unit]
Description=Atlas Control — Offline GPS Nav + Mesh Communications
Wants=network-online.target bluetooth.service ollama.service
After=network-online.target bluetooth.service nginx.service ollama.service

[Service]
Type=simple
User=ATLAS_USER
WorkingDirectory=APP_DIR
ExecStart=APP_DIR/venv/bin/python3 APP_DIR/app.py --host 0.0.0.0
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=ATLAS_UPS_I2C_BUS=7
Environment=ATLAS_UPS_I2C_ADDR=0x41

[Install]
WantedBy=multi-user.target
