# ollama.service
# NOTE: This file is a reference template only.
# The actual service file is generated by install.sh with the correct
# user and Ollama data directory substituted at install time.
# Do not copy this file directly to /etc/systemd/system/.

[Unit]
Description=Ollama LLM Service (Atlas Control)
After=network.target
Wants=network.target

[Service]
Type=simple
User=OLLAMA_USER
Group=OLLAMA_GROUP
Environment="OLLAMA_KEEP_ALIVE=10h"
Environment="OLLAMA_NUM_PARALLEL=1"
Environment="OLLAMA_GPU_OVERHEAD=0"
Environment="OLLAMA_MODELS=OLLAMA_DATA_DIR/models"
ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=5
TimeoutStartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ollama

[Install]
WantedBy=multi-user.target
